# fullmon.service
[Unit]
Description=fullmon

[Service]
ExecStart=/opt/fullmon/fullmon.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target