# fullmon.service
[Unit]
Description=fullmon

[Service]
ExecStart=/usr/local/bin/fullmon.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target